# systemd/nodejs.socket
[Unit]
Description=Socket for NodeJS Systemd Issue
PartOf=nodejs.service

[Socket]
ListenStream=8020

[Install]
WantedBy=sockets.target
